# roles/web_app/templates/web_app.socket.j2
{{ ansible_managed | comment }}

[Unit]
Description=sockets for the {{ web_app_name }} application

[Socket]
ListenStream=80
ListenStream=443

Backlog=128
NoDelay=true
ReusePort=true
