# roles/seahorse/templates/seahorse.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=seahorse
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/seahorse

[Install]
WantedBy=multi-user.target
