# roles/go/templates/go.socket.j2
[Unit]
Description=sockets for the {{ go_app_name }} application

[Socket]
ListenStream=80
ListenStream=443
NoDelay=true
